# dotfm/deploy/etc/s6-overlay/s6-rc.d/server/finish
#!/command/execlineb -S0

foreground { redirfd -w 1 /run/s6-linux-init-container-results/exitcode echo 0 }
/run/s6/basedir/bin/halt